# com.coretex.aos-organ.soul-dreamer.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.coretex.aos-organ.soul-dreamer</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>scripts/run-dream.js</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Library/AI/AI-AOS/AOS-organ-dev/AOS-organ-soul/AOS-organ-soul-src</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>DREAM_ENABLED</key>
        <string>false</string>
        <key>NODE_ENV</key>
        <string>development</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Library/AI/AI-Infra-MDvaults/MDvault-LLM-Ops/210-Log/soul-logs/soul-dreamer.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Library/AI/AI-Infra-MDvaults/MDvault-LLM-Ops/210-Log/soul-logs/soul-dreamer.stderr.log</string>
</dict>
</plist>
